# update.cron
0 * * * * root /usr/local/pf-socks5/update.sh >> /var/log/pf-socks5.log
